# autostart/yle-autostart.desktop
[Desktop Entry]
Name=Yle Autostart
Exec=/home/yle/Desktop/run.sh
Type=Application
Terminal=false
